# .config/autostart/Set caps as ctrl key.desktop
[Desktop Entry]
Exec=setxkbmap -option ctrl:swapcaps
Name=Set caps as ctrl key
OnlyShowIn=LXQt;
Type=Application
Version=1.0
